# tomatovirus_Inception V3/supervisor.conf
[program:tomato_virus]
command=/home/micz747/anaconda3/envs/tomatovirus/bin/python application.py
directory=/home/micz747/tomatovirus/
numprocs=1
autostart=true
startsecs=30
user=serving
